# systemd/pktgate.service
[Unit]
Description=pktgate - eBPF XDP Packet Filter
Documentation=file:///opt/pktgate/ARCHITECTURE.md
After=network-online.target
Wants=network-online.target

[Service]
# Type=notify enables READY=1 ack and the WatchdogSec= ping cycle. pktgate_ctl
# writes to $NOTIFY_SOCKET directly (no libsystemd dep).
Type=notify
NotifyAccess=main
EnvironmentFile=-/etc/pktgate/pktgate.conf
ExecStart=/bin/sh -c '/usr/local/bin/pktgate_ctl --json ${PKTGATE_METRICS_PORT:+--metrics-port $PKTGATE_METRICS_PORT} ${PKTGATE_EXTRA_ARGS} ${PKTGATE_CONFIG:-/etc/pktgate/config.json}'
ExecReload=/bin/kill -HUP $MAINPID

# Graceful shutdown — prints final stats before exit
KillSignal=SIGTERM
TimeoutStopSec=5

# Watchdog: kill+restart if the 1s main loop misses 30 pings in a row.
# Crash with XDP loaded → kernel keeps the last good program attached and
# traffic keeps flowing; stop → loader.detach() unhooks XDP → traffic passes
# unfiltered. Document this in ARCHITECTURE.md §Fail-safe.
WatchdogSec=30

# Restart policy
Restart=on-failure
RestartSec=3

# Security hardening
# AmbientCapabilities are inherited at exec from the unit; NoNewPrivileges=yes
# does NOT block them (that flag only stops new privileges via setuid/setgid,
# capabilities granted by the unit pass through).
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/sys/fs/bpf
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
RestrictNamespaces=yes
SystemCallArchitectures=native
MemoryDenyWriteExecute=yes

# Narrow the syscall surface to what an eBPF/network daemon needs.
# Cuts post-exploit blast radius — even if capabilities are subverted, the
# attacker cannot reach mount(), ptrace(), kexec_load(), etc.
SystemCallFilter=@system-service @network-io @file-system
SystemCallFilter=~@mount @reboot @swap @raw-io @cpu-emulation @obsolete
SystemCallErrorNumber=EPERM

# Resource limits — BPF maps require locked memory
LimitMEMLOCK=unlimited

# BPF + XDP attach on kernels >= 5.8 need CAP_BPF + CAP_NET_ADMIN + CAP_PERFMON.
# CAP_SYS_ADMIN is intentionally absent — it would grant mount/ns/kexec/raw-io
# rights orthogonal to BPF, defeating the rest of this hardening block.
AmbientCapabilities=CAP_BPF CAP_NET_ADMIN CAP_PERFMON
CapabilityBoundingSet=CAP_BPF CAP_NET_ADMIN CAP_PERFMON

# Logging to journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pktgate_ctl

[Install]
WantedBy=multi-user.target
