# systemd/pktgate.service
[Unit]
Description=pktgate - eBPF XDP Packet Filter
Documentation=file:///opt/pktgate/ARCHITECTURE.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/pktgate/pktgate.conf
ExecStart=/bin/sh -c '/usr/local/bin/pktgate_ctl --json ${PKTGATE_METRICS_PORT:+--metrics-port $PKTGATE_METRICS_PORT} ${PKTGATE_EXTRA_ARGS} ${PKTGATE_CONFIG:-/etc/pktgate/config.json}'
ExecReload=/bin/kill -HUP $MAINPID

# Graceful shutdown — prints final stats before exit
KillSignal=SIGTERM
TimeoutStopSec=5

# Restart policy
Restart=on-failure
RestartSec=3

# Security hardening
NoNewPrivileges=no
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/sys/fs/bpf
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
RestrictNamespaces=yes
SystemCallArchitectures=native
MemoryDenyWriteExecute=yes

# Resource limits — BPF maps require locked memory
LimitMEMLOCK=unlimited

# BPF requires these capabilities
AmbientCapabilities=CAP_BPF CAP_NET_ADMIN CAP_SYS_ADMIN CAP_PERFMON
CapabilityBoundingSet=CAP_BPF CAP_NET_ADMIN CAP_SYS_ADMIN CAP_PERFMON

# Logging to journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pktgate_ctl

[Install]
WantedBy=multi-user.target
